# infrastructure/systemd/docker.service.d/startup-delay.conf
# Delay dockerd (and all containers) on boot so the guest OS fully settles
# before the container I/O storm — added after the VM-100 NVMe-under-load
# incident (link dropped off the bus during a backup + mass-search storm).
#
# This is the unclean-boot fallback: after a power loss, ExecStop in
# homestack-startup.service didn't run, so the daemon still auto-restarts every
# restart:unless-stopped container at once — but only after the OS/NFS have had
# 90s to settle. On a CLEAN reboot the wave startup (homestack-startup.service)
# does the ordered bring-up instead.
#
# Install: /etc/systemd/system/docker.service.d/startup-delay.conf
# (handled by `hs startup` / scripts/staggered-up.sh --install)
[Service]
ExecStartPre=/bin/sleep 90
