# infrastructure/systemd/homestack-startup.service
[Unit]
Description=Homestack staggered startup — bring the stack up in ordered waves (avoids the boot I/O storm)
Documentation=file:///opt/docker/stacks/scripts/staggered-up.sh
# Start after the daemon, network, and remote mounts are up. The docker
# startup-delay drop-in (docker.service.d/startup-delay.conf) already holds
# dockerd 90s on boot so the guest OS/NFS settle before any container starts.
Wants=docker.service network-online.target
After=docker.service network-online.target remote-fs.target
OnFailure=homestack-notify-failure@%n.service

[Service]
Type=oneshot
RemainAfterExit=yes
# ExecStart brings every enabled stack/service up in waves with a settle gap.
ExecStart=/opt/docker/stacks/scripts/staggered-up.sh --yes
# On shutdown/reboot, mark every container user-stopped so the daemon does NOT
# auto-restart them concurrently on the next boot — this unit becomes the sole,
# ordered starter. (Runs while docker is still up, before docker.service stops.)
ExecStop=/opt/docker/stacks/scripts/staggered-up.sh --stop --yes
# Be gentle and patient — the waves + settle gaps take a few minutes.
Nice=5
TimeoutStartSec=20min
TimeoutStopSec=5min

[Install]
WantedBy=multi-user.target
