# infrastructure/systemd/homestack-backup.timer
[Unit]
Description=Run the homestack backup nightly at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
# Catch up a missed run after the box was off/rebooting at 03:00.
Persistent=true
# Spread I/O if other timers fire at the same minute.
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
