# infrastructure/systemd/homestack-backup.service
[Unit]
Description=Homestack backup — DB dumps + config/photo/document mirror to the backup share
Documentation=file:///opt/docker/stacks/scripts/backup.sh
Wants=docker.service network-online.target
After=docker.service network-online.target remote-fs.target
# NOTE: no RequiresMountsFor for the NFS share. It made a not-yet-ready mount
# after a reboot a hard dependency failure that silently abandoned the night's
# run. backup.sh now waits for its mounts (wait_mount) and aborts loudly if the
# configured target isn't really mounted, so OnFailure still fires on real loss.
OnFailure=homestack-notify-failure@%n.service

[Service]
Type=oneshot
ExecStart=/opt/docker/stacks/scripts/backup.sh
# Be gentle — this runs while other services are live.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=2h
