# infrastructure/systemd/homestack-notify-failure@.service
[Unit]
Description=ntfy alert for failed unit %i

[Service]
Type=oneshot
ExecStart=/opt/docker/stacks/scripts/notify-failure.sh %i
